# crontab
# Run every one minute
* * * * * /root/cloudflare-ddns.sh